# etc/systemd/system/baryssh-client.service
[Unit]
Description=barySSH client service

[Service]
Type=oneshot
EnvironmentFile=/etc/baryssh-client.conf
ExecStart=baryssh
