# etc/systemd/system/baryssh-server.service
[Unit]
Description=barySSH server service

[Service]
Type=oneshot
EnvironmentFile=/etc/baryssh-server.conf
ExecStart=baryssh
